# microservices/weather/scripts/weather.service
[Unit]
Description=Weather Service

[Service]
ExecStart=/usr/bin/weather-service

[Install]
WantedBy=multi-user.target
